# BiliBiliHelper.conf
[program:BiliBiliHelper]
process_name=%(program_name)s
command=python3 /root/BiliBiliHelper/main.py -d
autostart=true
autorestart=true
redirect_stderr=true